# projects/kairos/backend/deployments/systemd/icon-downloader.service
[Unit]
Description=ICON-EU Weather Data Downloader
After=network.target postgresql.service

[Service]
Type=oneshot
User=sellinios
Group=sellinios
WorkingDirectory=/home/sellinios/development/projects/aethra/backend
Environment="PATH=/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DB_PASSWORD=kairos123"
ExecStart=/home/sellinios/development/projects/aethra/backend/scripts/icon_update.sh
StandardOutput=journal
StandardError=journal

# Restart on failure with delay
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target